# bcm_userspace/pkg/files/etc/init.d/fap
#!/bin/sh /etc/rc.common
#set -x
START=21

start() {
     setupfap wan
     setupfap lan
}
