# buttons/pkg/files/etc/hotplug.d/button/buttons
#!/bin/sh

logger $ACTION
logger $BUTTON
